# etc/supervisor/conf.d/dante-socks.conf
[program:dante-socks]
command=/opt/dante-socks/sbin/sockd -f /etc/dante-socks/sockd.conf
user=root
autostart=true
autorestart=true
startsecs=10
